# services/apache2/run
#!/bin/sh
set -e
. /etc/apache2/envvars
exec /usr/sbin/apache2 -DFOREGROUND
