# services/phd/run
#!/bin/sh
set -e
exec /opt/phabricator/bin/phd start
